# systemd/llama-server.service
[Unit]
Description=llama.cpp local LLM server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/hyukju/work/oracle
EnvironmentFile=-/home/hyukju/work/oracle/.env
Environment=LLAMA_MODEL=/home/hyukju/work/oracle/models/gemma-4-E2B-it-UD-Q2_K_XL.gguf
ExecStart=/home/hyukju/work/oracle/scripts/run_llama_server.sh
Restart=on-failure
RestartSec=5
User=hyukju

[Install]
WantedBy=multi-user.target
