# systemd/llama-server.service
[Unit]
Description=llama.cpp local LLM server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/willtek/work/oracle
Environment=LLAMA_MODEL=/home/willtek/work/oracle/models/model.gguf
ExecStart=/home/willtek/work/oracle/scripts/run_llama_server.sh ${LLAMA_MODEL}
Restart=on-failure
RestartSec=5
User=willtek

[Install]
WantedBy=multi-user.target
