# systemd/oracle-report.service
[Unit]
Description=Oracle on-device Flask app
After=network-online.target llama-server.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/tmdwls/projects/oracle
EnvironmentFile=-/home/tmdwls/projects/oracle/.env
ExecStart=/home/tmdwls/projects/oracle/run.sh
Restart=on-failure
RestartSec=5
User=tmdwls

[Install]
WantedBy=multi-user.target
